# templates/default/thumbor-init-redhat.erb
#!/bin/bash

# WARNING: This file is managed by Chef, any manual changes will be OVERWRITTEN.

### BEGIN INIT INFO
# Provides:          thumbor
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Thumbor image manIPulation service
# Description:       Thumbor is a Python service providing an HTTP GET API to manipulate images
#                    This startup script launch as much as service as ports defined in /etc/default/thumbor
### END INIT INFO


PATH=/sbin:/usr/sbin:/bin:/usr/:/usr/local/bin:/usr/bin
NAME=thumbor

CONFFILE=/etc/thumbor.conf
KEYFILE=/etc/thumbor.key
IP=0.0.0.0
PORT=9000

DAEMON=`which thumbor`
DAEMON_ARGS="-c ${CONFFILE} -i ${IP} -k ${KEYFILE}"
PIDDIR=/var/run/$NAME
USER=thumbor
LOCKFILE=/var/lock/subsys/thumbor

# Source function library.
INITD=/etc/rc.d/init.d
. $INITD/functions

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

RETVAL=0
#
# Function that starts the daemon/service
#
start()
{
  echo -n "Starting $NAME: "
  touch /var/log/$NAME.log
  chown $USER /var/log/$NAME.log
  for p in `echo $PORT | tr ',' ' '`; do
    daemon --user $USER --pidfile $PIDDIR/$p.pid "nohup $DAEMON $DAEMON_ARGS -p ${p} 2>&1 1>/var/log/${NAME}-${p}.log &"
  done
  RETVAL=$?
  echo
  [ $RETVAL = 0 ] && touch $LOCKFILE
  return $RETVAL
}

#
# Function that stops the daemon/service
#
stop()
{
  echo -n "Stopping $NAME: "
  killproc $NAME
  RETVAL=$?
  echo
  [ $RETVAL = 0 ] && rm -f $LOCKFILE
  return $RETVAL
}

case "$1" in
  start)
    start
  ;;
  stop)
    stop
  ;;
  status)
    status $NAME
  ;;
  restart|force-reload)
    stop
    start
  ;;
  *)
  echo "Usage: $SCRIPTNAME {start|stop|status|restart|force-reload}" >&2
  exit 3
  ;;
esac

exit $RETVAL

# vim: ft=sh
